# Alone/dashboard-cron-conf/crontab
* * * * * root /var/www/html/dashboardSmartCity/ScheduledTasks/run-iotapp.sh > /proc/1/fd/1 2>/proc/1/fd/2
